# deluged.conf
[program:deluged]
autorestart=true
command=/run_deluge.sh

