# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:deluged]
autorestart=true
command=/run_deluge.sh

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:vnstatd]
command=/bin/su -c '/usr/sbin/vnstatd -n' deluge
autorestart=true

