# vnstatd.conf
[program:vnstatd]
command=/bin/su -c '/usr/sbin/vnstatd -n' deluge
autorestart=true

